# library/old-landchad/calibre-server.service
[Unit]
Description=Calibre library server
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/calibre-server --enable-auth --enable-local-write /opt/calibre/andersons --listen-on 127.0.0.1

[Install]
WantedBy=multi-user.target